# systemd/litellm.service
[Unit]
Description=LiteLLM gateway for CodingCrew
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/CodingCrew
EnvironmentFile=%h/CodingCrew/.env
ExecStart=/bin/bash -lc 'litellm --config %h/CodingCrew/config/litellm.yaml --port 4000 --host 127.0.0.1'
Restart=always
RestartSec=5
StandardOutput=append:%h/CodingCrew/logs/litellm.log
StandardError=append:%h/CodingCrew/logs/litellm.log

[Install]
WantedBy=default.target
