# systemd/orchestrator.service
[Unit]
Description=CodingCrew orchestrator
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
Environment="HOME=__HOME_DIR__"
Environment="XDG_CONFIG_HOME=__HOME_DIR__/.config"
EnvironmentFile=__REPO_DIR__/.env
ExecStartPre=/bin/bash -lc 'source __REPO_DIR__/.env && echo "$GITHUB_TOKEN" | gh auth login --with-token 2>/dev/null || true'
ExecStart=__REPO_DIR__/.venv/bin/python3 -m src.orchestrator
Restart=on-failure
RestartSec=30
RuntimeMaxSec=21600
StandardOutput=append:__REPO_DIR__/logs/orchestrator.log
StandardError=append:__REPO_DIR__/logs/orchestrator.log

[Install]
WantedBy=default.target
